# deployment/systemd/hodlxxi-mcp.service
[Unit]
Description=HODLXXI read-only MCP localhost sidecar
Documentation=https://github.com/hodlxxi/Universal-Bitcoin-Identity-Layer
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
DynamicUser=yes
User=hodlxxi-mcp
Group=hodlxxi-mcp

WorkingDirectory=/opt/hodlxxi-mcp/current
ExecStart=/opt/hodlxxi-mcp/current/venv/bin/hodlxxi-mcp-http

Restart=on-failure
RestartSec=5s
TimeoutStartSec=30s
TimeoutStopSec=10s
KillSignal=SIGINT
KillMode=mixed

Environment=HOME=/tmp
Environment=PYTHONDONTWRITEBYTECODE=1
Environment=PYTHONUNBUFFERED=1
UMask=0077

NoNewPrivileges=yes
PrivateTmp=yes
PrivateDevices=yes
DevicePolicy=closed
ProtectSystem=strict
ProtectHome=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
ProtectHostname=yes
ProtectClock=yes
ProtectProc=invisible
ProcSubset=pid
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
RemoveIPC=yes
KeyringMode=private
SystemCallArchitectures=native
CapabilityBoundingSet=
AmbientCapabilities=

SocketBindDeny=any
SocketBindAllow=tcp:8765

InaccessiblePaths=-/srv/ubid
InaccessiblePaths=-/srv/ubid-staging
InaccessiblePaths=-/etc/hodlxxi
InaccessiblePaths=-/var/lib/lnd
InaccessiblePaths=-/home/lnd
ReadOnlyPaths=-/opt/hodlxxi-mcp

LimitNOFILE=4096
TasksMax=128
MemoryMax=512M
OOMPolicy=stop

StandardOutput=journal
StandardError=journal
SyslogIdentifier=hodlxxi-mcp

[Install]
WantedBy=multi-user.target
